# debian/unison-cronjob.cron.d
*/15 *  *   *   *   user  /usr/sbin/unison-cronjob
